# systemd/bsc-exporter.service
[Unit]
Description=BSC → AWS Public Blockchain daily export
Documentation=file:///opt/bsc-exporter/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot

# Adjust these to match your deployment.
User=bsc
Group=bsc
WorkingDirectory=/opt/bsc-exporter

# Environment for AWS credentials. Prefer IAM role on EC2 — leave these
# unset and boto3 will pick up the instance profile automatically.
# Otherwise drop credentials into /etc/bsc-exporter/aws.env (chmod 600)
# and uncomment:
# EnvironmentFile=/etc/bsc-exporter/aws.env

# Run against the last 7 days so transient failures (RPC blip, node OOM,
# etc.) get auto-corrected within a week — progress.txt + S3 skip make all
# already-done dates near-free, so the actual work each day is usually just
# the new "yesterday". Bash wrapper is needed for $(date) substitution;
# %%F is systemd's escape for a literal % in ExecStart.
# Add `--collector.textfile.directory=/var/lib/node_exporter/textfile_collector`
# to PATH if you wire textfile metrics later.
Environment=PATH=/opt/bsc-exporter/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/bin/bash -c '/opt/bsc-exporter/venv/bin/python /opt/bsc-exporter/exporter.py --config /etc/bsc-exporter/config.yaml --start "$(date -u -d "7 days ago" +%%F)" --end "$(date -u -d "yesterday" +%%F)"'

# Hardening — safe defaults for a network/disk-bound script.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ReadWritePaths=/var/lib/bsc-exporter /var/log/bsc-exporter

# Logs go to journald — view with `journalctl -u bsc-exporter`.
StandardOutput=journal
StandardError=journal

# A single export run can take 30+ min on slow RPCs; give it room.
TimeoutStartSec=2h

[Install]
WantedBy=multi-user.target
